# deploy/ssh-drop-restart.service
[Unit]
Description=Restart ssh-drop server

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl restart ssh-drop.service
