# deploy/ssh-drop-restart.timer
[Unit]
Description=Daily restart of ssh-drop server

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
